# lib/generators/capistrano/unicorn_nginx_osx/templates/unicorn.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>EnvironmentVariables</key>
  <dict>
      <key>PATH</key>
      <string>/Users/<%= fetch(:deploy_user)%>/.rbenv/shims:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>RBENV_ROOT</key>
      <string>/Users/<%= fetch(:deploy_user)%>/.rbenv</string>
      <key>RBENV_VERSION</key>
      <string>2.2.3</string>
  </dict>
  <key>Label</key>
  <string>apps.<%= fetch(:application)%>.unicorn</string>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
        <string>/Users/<%= fetch(:deploy_user)%>/.rbenv/shims/unicorn</string>
        <string>-c</string>
        <string>/Users/<%= fetch(:deploy_user)%>/apps/<%= fetch(:application) %>/shared/config/unicorn.rb</string>
        <string>-E</string>
        <string>production</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>